# wordpress/services.d/caddy/run
#!/usr/bin/execlineb -P

with-contenv
s6-envuidgid www-data

exec caddy -agree -conf	 /Caddyfile -email basal@${ROOT_DOMAIN}
